# profile.d/jboss_profile.sh
export JBOSS_HOME="/opt/JBOSS_UNZIP_PATH/jboss-eap-7.2"
export JBOSS_BIN="/opt/JBOSS_UNZIP_PATH/jboss-eap-7.2/bin"
